# ec2-api/deploy/systemd/hosted-skill-api.service
[Unit]
Description=Hosted skill API (FastAPI / uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/hosted-skill-api/ec2-api
EnvironmentFile=/opt/hosted-skill-api/ec2-api/.env
# .venv must be created with Python 3.14.x (see README.md / pyproject.toml)
ExecStart=/opt/hosted-skill-api/ec2-api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
